# graphical-session.target
[Unit]
Description=X session is up (started from .xinitrc)
Requires=dbus.service
After=dbus.service
